# roles/zookeeper/templates/supervisor-zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper-{{zookeeper.version}}/bin/zkServer.sh start-foreground
directory=/opt/zookeeper-{{zookeeper.version}}
redirect_stderr=true
user=zookeeper
