# roles/zookeeper/files/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper-3.4.5/bin/zkServer.sh start-foreground
directory=/opt/zookeeper-3.4.5
redirect_stderr=true
